# deploy/com.juggle.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.juggle.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>-S</string>
        <string>uv</string>
        <string>run</string>
        <string>python</string>
        <string>src/juggle_watchdog_daemon.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mikechen/github/juggle</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>JUGGLE_ORCHESTRATOR</key>
        <string>1</string>
        <key>JUGGLE_WATCHDOG_SUPERVISED</key>
        <string>1</string>
        <key>PATH</key>
        <string>/Users/mikechen/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/mikechen/.juggle/watchdog.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mikechen/.juggle/watchdog.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
